# deploy/cryptoforge.service
# ═══════════════════════════════════════════════════════════════
#  CryptoForge — systemd template unit for blue-green deploys
#  Install as: /etc/systemd/system/cryptoforge@.service
#  Usage:      systemctl start cryptoforge@9000
#              systemctl start cryptoforge@9001
#  The %i expands to the port number after the @ sign.
#  Port is passed directly via --port %i; APP_PORT in .env is ignored for
#  the bind address to avoid load_dotenv() overriding the systemd value.
# ═══════════════════════════════════════════════════════════════
[Unit]
Description=CryptoForge Crypto Trading Platform (port %i)
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/cryptoforge
EnvironmentFile=/home/ec2-user/cryptoforge/.env
ExecStart=/home/ec2-user/cryptoforge/venv/bin/uvicorn app:app --host 127.0.0.1 --port %i --log-level info
Restart=on-failure
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGINT
KillMode=mixed
UMask=0077
LimitNOFILE=65536

# ── Memory ceiling ──────────────────────────────────────────────
# Two of these run at once during a blue-green deploy, on a small Lightsail
# box that also hosts nginx. Without a ceiling a runaway process is reaped by
# the kernel OOM killer, which picks its own victim — and that has taken nginx
# down with it, so the failure reads as "the whole site is gone" rather than
# "one service restarted". With MemoryMax set, the process dies alone and
# Restart=on-failure brings it straight back.
# Raise both if the instance is larger than 1 GB; MemoryHigh should stay a
# little under MemoryMax so reclaim gets a chance before the hard kill.
MemoryAccounting=true
MemoryHigh=280M
MemoryMax=340M
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
LockPersonality=true
ProcSubset=pid
ProtectProc=invisible
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
