# deploy/cryptoforge.service
# ═══════════════════════════════════════════════════════════════
#  CryptoForge — systemd template unit for blue-green deploys
#  Install as: /etc/systemd/system/cryptoforge@.service
#  Usage:      systemctl start cryptoforge@9000
#              systemctl start cryptoforge@9001
#  The %i expands to the port number after the @ sign.
#  APP_PORT env var overrides config.py's default.
# ═══════════════════════════════════════════════════════════════
[Unit]
Description=CryptoForge Crypto Trading Platform (port %i)
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/cryptoforge
EnvironmentFile=/home/ec2-user/cryptoforge/.env
Environment="APP_PORT=%i"
ExecStart=/home/ec2-user/cryptoforge/venv/bin/python3 app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
